# meta-txgrl/recipes-core/systemd/files/init_cts3.service
[Unit]
Description=Init CTS3 board.

[Service]
Type=oneshot
ExecStart=/usr/bin/init_cts3.sh

[Install]
WantedBy=multi-user.target
