# meta-txgrl/recipes-core/systemd/files/rngd.service
[Service]
ExecStart=/usr/sbin/rngd -f -r /dev/urandom
